# home/.config/systemd/user/omniroute.service
[Unit]
Description=omniroute

[Service]
ExecStart=/home/rey/.local/npm-global/bin/omniroute
Restart=always

[Install]
WantedBy=multi-user.target
